# blog-service/src/debian/data/etc/init/blog-service.conf
description "Dropwizard blog service"
author "Charlie Knudsen <charlie.knudsen@gmail.com>"

start on (filesystem and net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 5 30

script
  export DW_DB_URL="jdbc:postgresql://127.0.0.1:5432/blog"
  export DW_DB_PASSWORD="blog"
  export ETCD_URL="http://10.0.2.15:4001"
  export HOST_IP=`/usr/bin/iplookup 'eth0 '`
  export APP_OPTS="-Ddw.database.url=$DW_DB_URL -Ddw.database.password=$DW_DB_PASSWORD -Ddw.etcd.hosts=$ETCD_URL -Ddw.etcd.publish.hostName=$HOST_IP"
  export APP_JAR="/usr/share/blog-service/application.jar"

  export GC_OPTS="-XX:+UseParallelOldGC"
  export JAVA_OPTS="-server -Xms256M -Xmx512M"

  /usr/bin/java $JAVA_OPTS $GC_OPTS $APP_OPTS -jar $APP_JAR server /usr/share/blog-service/config.yaml
end script
